# .config/autostart/Brightness.desktop
[Desktop Entry]
Type=Application
Exec="/home/subin/conf/brightness.sh"
StartupNotify=false
